# server/src/main/rpm/scriptlets/post-install.sh
#!/bin/sh

/sbin/chkconfig --add office-account-manager
